# netflix-stack-poc/dockerfiles/hystrix-dashboard/files/supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslogd]
command=/usr/sbin/rsyslogd -n

[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:tomcat7]
command=/usr/share/tomcat7/bin/catalina.sh run
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
environment=CATALINA_BASE="/var/lib/tomcat7"